# immortalwrt/build.sh
#!/bin/bash
set -euxo pipefail

# ===============================
# OneCloud ImmortalWrt Build Script
# ===============================

echo "ðŸš€ å¼€å§‹å‡†å¤‡æž„å»ºçŽ¯å¢ƒ..."

# æ›´æ–° feeds ç¡®ä¿ä¾èµ–å®Œæ•´
./scripts/feeds update -a
./scripts/feeds install -a

# æ¸…ç†ç¼“å­˜é¿å…ä¸Šæ¬¡æ®‹ç•™é”™è¯¯
make clean || true
rm -rf tmp/ || true

# ===============================
# è‡ªå®šä¹‰è¦å®‰è£…çš„åŒ…
# ===============================
PACKAGES=""

# åŸºç¡€åŠŸèƒ½ & ä¸­æ–‡ç•Œé¢
PACKAGES="$PACKAGES curl wget ca-certificates"
PACKAGES="$PACKAGES luci luci-compat luci-base luci-app-firewall"
PACKAGES="$PACKAGES luci-i18n-base-zh-cn luci-i18n-firewall-zh-cn luci-i18n-package-manager-zh-cn"
#PACKAGES="$PACKAGES ppp ppp-mod-pppoe luci-proto-ppp"

# å¸¸ç”¨å¢žå¼ºæ’ä»¶
PACKAGES="$PACKAGES luci-app-docker luci-app-ttyd luci-app-filebrowser"
PACKAGES="$PACKAGES kmod-usb-storage block-mount e2fsprogs fdisk"
PACKAGES="$PACKAGES luci-app-opkg openssh-sftp-server"

# ===============================
# å†™å…¥æ‰©å®¹è„šæœ¬ï¼ˆç³»ç»Ÿå¯åŠ¨è‡ªåŠ¨æ‰§è¡Œï¼‰
# ===============================
echo "ðŸ§© æ·»åŠ è‡ªåŠ¨æ‰©å®¹è„šæœ¬..."
mkdir -p files/etc/init.d

cat > files/etc/init.d/expand_rootfs <<'EOF'
#!/bin/sh /etc/rc.common
START=99
DESCRIPTION="Auto expand root filesystem on first boot"

start() {
    if [ ! -f /etc/expand_done ]; then
        echo "ðŸ”§ æ­£åœ¨è‡ªåŠ¨æ‰©å±• eMMC åˆ†åŒº..."

        parted /dev/mmcblk1 resizepart 2 100%
        losetup /dev/loop0 /dev/mmcblk1p2
        e2fsck -f -y /dev/loop0
        resize2fs -f /dev/loop0
        sync

        echo "âœ… åˆ†åŒºæ‰©å±•å®Œæˆã€‚ç³»ç»Ÿå°†è‡ªåŠ¨é‡å¯ä»¥ç”Ÿæ•ˆ..."
        touch /etc/expand_done
        reboot
    fi
}
EOF

chmod +x files/etc/init.d/expand_rootfs

# ===============================
# æž„å»ºé•œåƒï¼ˆæ ¹åˆ†åŒºè°ƒæ•´ä¸º 1024MBï¼‰
# ===============================
echo "ðŸ§± å¼€å§‹æž„å»ºé•œåƒ..."
make image \
  PACKAGES="$PACKAGES" \
  FILES="files" \
  ROOTFS_PARTSIZE="1024" \
  V=s

# ===============================
# è‡ªåŠ¨ç”Ÿæˆæ›´æ–°è¯´æ˜Žï¼ˆä¾› Release ä½¿ç”¨ï¼‰
# ===============================
echo "ðŸ“„ ç”Ÿæˆæ›´æ–°è¯´æ˜Ž..."
mkdir -p ../release_note
cat > ../release_note/update.txt <<EOF
ðŸ†• æœ¬æ¬¡æž„å»ºæ›´æ–°å†…å®¹ï¼š
- è‡ªåŠ¨æ‰©å±• eMMC åˆ†åŒºï¼ˆé¦–æ¬¡å¯åŠ¨è‡ªåŠ¨å®Œæˆï¼‰
- ä¿®å¤ ppp-mod-pppoe å®‰è£…å¤±è´¥é—®é¢˜
- æ–°å¢ž luci-app-dockerã€luci-app-ttydã€luci-app-filebrowser ç­‰æ’ä»¶
- å®Œå–„ä¸­æ–‡ç•Œé¢æ”¯æŒ
- é»˜è®¤ root å¯†ç ä¸ºç©ºï¼ˆç›´æŽ¥ç™»å½•ï¼‰
- æ ¹åˆ†åŒºæ‰©å¤§ä¸º 1024MBï¼Œç©ºé—´æ›´å……è¶³
EOF

echo "âœ… æž„å»ºå®Œæˆï¼å›ºä»¶è¯·æŸ¥çœ‹ bin/targets/"

